# roles/openclaw/templates/clawdbot-host.service.j2
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ clawdbot_home }}

# Environment variables
Environment="PATH={{ clawdbot_home }}/.local/share/pnpm:{{ clawdbot_home }}/.local/bin:{{ clawdbot_home }}/.npm-global/bin:{{ clawdbot_home }}/bin:{{ clawdbot_home }}/.nvm/current/bin:{{ clawdbot_home }}/.fnm/current/bin:{{ clawdbot_home }}/.volta/bin:{{ clawdbot_home }}/.asdf/shims:{{ clawdbot_home }}/.bun/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PNPM_HOME={{ clawdbot_home }}/.local/share/pnpm"

# Start command
{% if clawdbot_install_mode == 'release' -%}
ExecStart={{ clawdbot_home }}/.local/share/pnpm/openclaw gateway
{% else -%}
ExecStart={{ clawdbot_home }}/.local/bin/clawdbot gateway
{% endif -%}

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

[Install]
WantedBy=default.target
